# compose/opennms/env/supervisord.conf
[supervisord]
nodaemon=true

[program:opennms_start]
command=/opennms_start.sh
stdout_logfile=/dev/stdout
